# docker-sip/config/etc/supervisord/supervisord.conf
[supervisord]
nodaemon = true
user = root
logfile = syslog 

[program:intercom]
autostart=true
autorestart=true
priority = 10
command = /app/intercom/intercom.py

